# systemd/infra-maintenance.timer
[Unit]
Description=Run maintenance daily at 02:00

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
